# deploy/com.snapstack.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  launchd LaunchAgent template for the snapstack server (macOS).
  Reference only — prefer ./install-macos.sh, which fills the absolute path and
  loads the agent for you. To use this file by hand, replace __LAUNCHER__ with the
  absolute path to deploy/snapstack-start.sh, copy it to ~/Library/LaunchAgents/,
  then run:
      launchctl load ~/Library/LaunchAgents/com.snapstack.server.plist
  The launcher (snapstack-start.sh) pulls the latest code and ensures deps before
  starting node, so the agent self-updates on each (re)start.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.snapstack.server</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/sh</string>
    <string>__LAUNCHER__</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/snapstack.out.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/snapstack.err.log</string>
  <!--
  <key>EnvironmentVariables</key>
  <dict>
    <key>SNAPSTACK_DIR</key><string>/Users/USER/.snapstack</string>
    <key>SNAPSTACK_PORT</key><string>4123</string>
  </dict>
  -->
</dict>
</plist>
